# supervisord.conf
[supervisord]
nodaemon=true

[program:autologin-server]
command=autologin-server --host 0.0.0.0

[program:autologin-http-api]
command=autologin-http-api
